# LINUX/perms/fix-permissions.service
[Unit]
Description=Fixes permissions issue caused by Windows

[Service]
ExecStart=/fix-permissions.sh

[Install]
WantedBy=multi-user.target
